# nginx/root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

# dns resolution fails if nginx is started too early
sleep 30;

exec nginx;
