# api/root/etc/services.d/mongodb/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid mongodb mongod --config /etc/mongod.conf;
